# conf/cron.d/suricata-rules-update
*/10 * * * * root suricata-update --reload-command "kill -usr2 $(pidof suricata)" 2>&1 >/var/log/suricata-rules-update.log
